# configs/supervisor-hermes.conf
[program:hermes]
command=/usr/bin/bash /opt/start-hermes-dashboard.sh
directory=/home/ubuntu
user=ubuntu
environment=HOME="/home/ubuntu",DISPLAY=":1",XDG_RUNTIME_DIR="/run/user/1000",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/hermes-agent"
priority=210
autostart=%(ENV_ENABLE_HERMES)s
autorestart=true
startretries=3
startsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/hermes_stdout.log
stderr_logfile=/tmp/hermes_stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
