# configs/supervisor-openclaw.conf
[program:openclaw]
command=/opt/start-openclaw.sh
directory=/home/ubuntu
user=ubuntu
environment=HOME="/home/ubuntu",DISPLAY=":1",XDG_RUNTIME_DIR="/run/user/1000"
priority=200
autostart=%(ENV_ENABLE_OPENCLAW)s
autorestart=false
startretries=5
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/openclaw_stdout.log
stderr_logfile=/tmp/openclaw_stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
