# configs/supervisor-dashboard.conf
[program:dashboard]
command=/opt/start-dashboard.sh
directory=/opt
user=ubuntu
environment=HOME="/home/ubuntu",NODE_ENV="production",NODE_PATH="/usr/local/lib/node_modules:/usr/lib/node_modules",AUTH_USER="%(ENV_AUTH_USER)s",AUTH_PASSWORD="%(ENV_AUTH_PASSWORD)s",OPENCLAW_GATEWAY_TOKEN="%(ENV_OPENCLAW_GATEWAY_TOKEN)s",VNC_PASSWORD="%(ENV_VNC_PASSWORD)s"
priority=100
autostart=true
autorestart=true
startretries=3
startsecs=5
stdout_logfile=/tmp/dashboard_stdout.log
stderr_logfile=/tmp/dashboard_stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
